# resources/infra/rc.local
#!/bin/sh -e

set -eu

export PORT=80
(cd /app; nodejs bin/server.js &)

exit 0
